# pn-vps/discovery-agent/privatenexus-discovery-agent.service
[Unit]
Description=PrivateNexus discovery agent (push host + container facts)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/privatenexus/scripts/discovery-agent.sh
